# thermal-manager.service
[Unit]
Description=Windows-like Thermal Management
After=multi-user.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/thermal-manager.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
